# scripts/ops/systemd/factory-monitor.timer
[Unit]
Description=Prismatic Factory Monitor — heartbeat health check (every 15 min)

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
